# bin/com.rondo.feedback-agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.rondo.feedback-agent</string>

	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>/Users/joost/Code/rondo/rondo-club/bin/get-feedback.sh</string>
		<string>--loop</string>
		<string>--optimize</string>
	</array>

	<key>StartInterval</key>
	<integer>300</integer>

	<key>WorkingDirectory</key>
	<string>/Users/joost/Code/rondo/rondo-club</string>

	<key>EnvironmentVariables</key>
	<dict>
		<key>HOME</key>
		<string>/Users/joost</string>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
	</dict>

	<key>StandardOutPath</key>
	<string>/Users/joost/Code/rondo/rondo-club/logs/launchd-stdout.log</string>

	<key>StandardErrorPath</key>
	<string>/Users/joost/Code/rondo/rondo-club/logs/launchd-stderr.log</string>

	<key>RunAtLoad</key>
	<false/>

	<key>KeepAlive</key>
	<false/>
</dict>
</plist>
